# files/kubelet-watchdog.service
[Unit]
Description=kubelet-watchdog

[Service]
User=root
Type=oneshot
ExecStart=/opt/bin/kubelet-watchdog
